# etc/s6-overlay/s6-rc.d/shairport/run
#!/command/with-contenv sh

echo "Starting Shairport Sync"
exec /usr/local/bin/shairport-sync $EXTRA_ARGS
